# gic/solitaire/rc.local
wine /root/wine/123freesolitaire-v1200-setup.exe
