# gic/csII/rc.local
wine /root/wine/Conflict_Desert_Storm_2/CDS\ II.EXE
